# server/trips-tg-brief.service
[Unit]
Description=Travel Companion daily Telegram briefing
After=network-online.target trips-sync.service

[Service]
Type=oneshot
EnvironmentFile=/opt/trips-sync/.env
# Hit the local trips-sync brief endpoint; key must equal TG_WEBHOOK_SECRET.
# (Adjust the port if uvicorn binds elsewhere — see trips-sync.service.)
ExecStart=/bin/sh -c 'curl -fsS "http://127.0.0.1:8090/tg/brief?key=$TG_WEBHOOK_SECRET"'
