# server/trips-sync.service
[Unit]
Description=Trips sync backend (alpine-loop-guide)
After=network.target

[Service]
WorkingDirectory=/opt/trips-sync
EnvironmentFile=/opt/trips-sync/.env
ExecStart=/opt/trips-sync/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8090
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target
